# rootfs/etc/s6-overlay/s6-rc.d/config-deno/run
#!/usr/bin/with-contenv sh

DENODIR=${DENO_DIR:-'/deno-dir'}

mkdir -p ${DENODIR}

setgroup ${DENODIR} /app
